# scripts/social-network-app@.service
[Unit]
Description=Social Network App — AI-assisted social media management dashboard
Documentation=https://github.com/techengineerworkstation/social-network-app
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=%i
Environment=NODE_ENV=production
Environment=PORT=8899
ExecStart=/usr/local/bin/social-network-app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=social-network-app

[Install]
WantedBy=multi-user.target
